# images/php/supervisor.conf
[supervisord]
user=root

[program:php-fpm]
command=php-fpm

[program:sshd]
command=/usr/sbin/sshd -D

[program:cron]
command=cron