# images/php/supervisor.conf
[supervisord]
user=root

[program:php-fpm]
command=php-fpm

[program:sshd]
command=/startup/sshd.sh

[program:cron]
command=cron